# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:backend]
command=java -Xmx384m -Djava.security.egd=file:/dev/./urandom -Dspring.datasource.url="${TIDB_URL:jdbc:mysql://your-tidb-host:4000/lawezy}" -Dspring.datasource.username="${TIDB_USER:your-user}" -Dspring.datasource.password="${TIDB_PASSWORD:your-password}" -Dspring.data.mongodb.uri="${MONGO_URI:mongodb+srv://user:pass@cluster.mongodb.net/lawezy}" -jar /app/backend/app.jar --server.port=8082 --server.address=0.0.0.0
directory=/app/backend
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:messenger]
command=node /app/messenger/server.js
directory=/app/messenger
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
